# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
user = root
command=uwsgi --socket apoio.sock --chdir /app/apoio --wsgi-file wsgi.py --chmod-socket=666
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
user = root
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
